# virtual/root/etc/systemd/system/catalis-poll@.service
[Unit]
Description=Catalis Cloud HMI MODBUS polling service
BindsTo=catalis-chmi@%i.service
Requires=catalis-init@%i.service
After=catalis-init@%i.service

[Service]
Type=simple
EnvironmentFile=/etc/catalis/global.conf
ExecStart=/usr/sbin/catalis-poll %i

[Install]
WantedBy=multi-user.target
